# etc/profile.d/00-pathfuncs.sh
# Path manipulation functions from BLFS
pathappend() {
    if [ -d "$1" ] && [[ ":$PATH:" != *":$1:"* ]]; then
        PATH="${PATH:+$PATH:}$1"
    fi
}

pathprepend() {
    if [ -d "$1" ] && [[ ":$PATH:" != *":$1:"* ]]; then
        PATH="$1${PATH:+:$PATH}"
    fi
}

export -f pathappend pathprepend
